# SOURCES/Leap-0.8.0/lib/systemd/system/leap.service
[Unit]
Description=LeapMotion Daemon
After=syslog.target

[Service]
Type=simple
ExecStart=/usr/bin/leapd

[Install]
WantedBy=multi-user.target